# uberspace/etc_services.d_uwsgi.ini
[program:uwsgi]
command=/home/statplan/.local/bin/uwsgi --master --emperor /home/<username>/uwsgi/apps-enabled
autostart=true
autorestart=true
stderr_logfile = ~/uwsgi/err.log
stdout_logfile = ~/uwsgi/out.log
stopsignal=INT